# memory-capper.service.in
[Unit]
Description=Memory Capper
Documentation=https://github.com/NickCis/memory-capper

[Service]
ExecStart=@bindir@/memory-capper --non-buffered-stdout --no-date --no-tick-log
# restart the process if it exits prematurely
Restart=on-failure
StartLimitBurst=3
StartLimitInterval=60s

[Install]
WantedBy=multi-user.target
